# deploy/pymash_monitoring.timer
[Unit]
Description=Runs monitoring every minute

[Timer]
OnUnitActiveSec=1m
Unit=pymash_monitoring.service

[Install]
WantedBy=multi-user.target
